# baza-task-runner.timer
[Unit]
Description=Baza Empire — Task Runner Timer (every 6 hours)
Requires=baza-task-runner.service

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
Unit=baza-task-runner.service

[Install]
WantedBy=timers.target
